# files/etc/rc.local
#!/bin/sh
# Enable Dropbear and set root password at boot

uci set dropbear.@dropbear[0].enable='1'
uci set dropbear.@dropbear[0].Port='22'
uci set dropbear.@dropbear[0].PasswordAuth='on'
uci set dropbear.@dropbear[0].RootPasswordAuth='on'
uci commit dropbear

/etc/init.d/dropbear enable
/etc/init.d/dropbear start

echo -e "root\nroot" | passwd root

exit 0
